# private_Library/LaunchAgents/com.joelazar.ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  User LaunchAgent: run `ollama serve` at login with performance env vars
  baked in. Replaces the previous setenv-only agent since we're using the
  Homebrew CLI (no Ollama.app installed).

  Why this over `brew services start ollama`?
    - brew services' plist has no env vars and is overwritten by brew upgrades.
    - This plist is fully owned by us, tracked in chezmoi, and survives
      `brew upgrade ollama` (binary path stays the same).

  Target hardware: MacBook Pro M4 Pro, 48 GB unified memory.
  Target model:    qwen3.6:35b-a3b (~24 GB) and similar MoE coding models.

  Manage with:
    launchctl load   ~/Library/LaunchAgents/com.joelazar.ollama.plist
    launchctl unload ~/Library/LaunchAgents/com.joelazar.ollama.plist
    launchctl kickstart -k gui/$(id -u)/com.joelazar.ollama    # restart

  Logs:
    tail -f /tmp/com.joelazar.ollama.out
    tail -f /tmp/com.joelazar.ollama.err
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joelazar.ollama</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ollama</string>
        <string>serve</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- 64K context — 32K fills up fast in an agent loop.
             KV cache at 64K q8_0 ≈ ~6 GB, leaves headroom on 48 GB. -->
        <key>OLLAMA_CONTEXT_LENGTH</key>
        <string>65536</string>

        <!-- Required for KV cache quantization. Pure win on Metal. -->
        <key>OLLAMA_FLASH_ATTENTION</key>
        <string>1</string>

        <!-- Halves KV cache RAM with near-zero quality loss. -->
        <key>OLLAMA_KV_CACHE_TYPE</key>
        <string>q8_0</string>

        <!-- Keep models resident for 30 min to avoid ~15 s cold reloads. -->
        <key>OLLAMA_KEEP_ALIVE</key>
        <string>30m</string>

        <!-- Single model, single parallel request — predictable RAM on 48 GB. -->
        <key>OLLAMA_NUM_PARALLEL</key>
        <string>1</string>

        <key>OLLAMA_MAX_LOADED_MODELS</key>
        <string>1</string>

        <!-- Ensure the binary finds Homebrew's env. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>

        <!-- Needed so HOME-relative paths (~/.ollama/models) resolve. -->
        <key>HOME</key>
        <string>/Users/joelazar</string>
    </dict>

    <!-- Start at login and keep alive (restart on crash). -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Throttle respawn if it keeps crashing. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/com.joelazar.ollama.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.joelazar.ollama.err</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
